# build/supervisor/etc/supervisor/conf.d/student_server.conf
[program:student_server]
# 运行的命令
command=/pitrix/lib/student-server/student_server -config /pitrix/config/stu_apiserver.yaml

# 表示是否跟随supervisor一起启动。
autostart=true

startretries=1000000

#如果该程序挂了，是否重新启动。
autorestart=true 

#终端标准输出重定向文件
#stdout_logfile

# 终端错误输出重定向文件
#stderr_logfile